# deployments/orion-agent-deployment/src/main/resources/deb/orion-agent.upstart
description "OrionAgent"

start on runlevel [2345]
setuid orion-agent
respawn
respawn limit 20 5

exec /opt/orion-agent/run_agent.sh